# Docker/supervisord/proftpd.conf
[program:proftpd]
command=proftpd --nodaemon
stdout_logfile=/var/log/supervisor/proftpd.log
stderr_logfile=/var/log/supervisor/proftpd.log
autorestart=false
startsecs=5
#user=root
